# mqtt_pigpio.service
[Unit]
Description=MQTT PIGPIO service
Requires=mosquitto.service
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=5
User=root
Environment="LD_LIBRARY_PATH=/usr/local/lib"
ExecStart=/usr/local/bin/mqtt_pigpio

[Install]
WantedBy=multi-user.target

